# beagleBone/valhalla.service.in
[Unit]
Description=Valhalla UDP Receiver

[Service]
WorkingDirectory=LEDSCAPE_PATH/
ExecStart=LEDSCAPE_PATH/run-valhalla
KillMode=control-group
Restart=always
Nice=-20

[Install]
WantedBy=multi-user.target
